# debian/DEBIAN/postinst
#!/bin/bash
set -e

# 在安装/升级后运行

APP_NAME="ems-simulate"
INSTALL_DIR="/usr/share/${APP_NAME}"
BIN_DIR="/usr/bin"
SERVICE_NAME="${APP_NAME}.service"

case "$1" in
    configure)
        # 1. 重新加载 systemd
        if command -v systemctl >/dev/null; then
            systemctl daemon-reload
            # 可选：安装后自动启动
            # systemctl enable "${SERVICE_NAME}"
            # systemctl start "${SERVICE_NAME}"
        fi
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0
